# heatmon.cron
CMD=/usr/lib/heatmon/heatmon.rb

* * * * * root $CMD
